# sparse/etc/systemd/system/systemd-udevd.service.d/nodevtmpfs.conf
[Service]
Environment="LD_PRELOAD=/usr/lib/systemd/udevd-nodevtmpfs.so"
